# runit/socklog-ucspi-tcp/run
#!/bin/sh
PORT=10116
exec 2>&1
exec tcpsvd -vllogserver -unobody 0 "$PORT" socklog ucspi TCPREMOTEIP
