# runit/getty-1/run
#!/bin/sh
exec /sbin/getty -8 38400 tty1
