# runit/example/run
#!/bin/sh

exec 2>&1
exec chpst -u nobody:nogroup -/ /srv/containers/base /bin/example
