# runit/nginx/run
#!/bin/sh

exec 2>&1
exec chpst   /usr/sbin/nginx -g "daemon off;"
